# config/block-announcer.service
[Unit]
Description=Schedule Announcer — AI voice + desktop notifications at schedule block changes

[Service]
Type=oneshot
ExecStart=/usr/bin/block-announcer
